# docs/com.luoqi.ccusage-sync.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd"?>
<!--
  launchd 定时同步 cc-usage 数据
  每 10 分钟跑一次 cc-usage sync，把 Claude + ZCode 日志同步到 ccusage.db
  （sandbox 下的 TokenMonitor App 无法 spawn 外部进程，所以走 launchd）

  安装：
    cp com.luoqi.ccusage-sync.plist ~/Library/LaunchAgents/
    launchctl load ~/Library/LaunchAgents/com.luoqi.ccusage-sync.plist

  卸载：
    launchctl unload ~/Library/LaunchAgents/com.luoqi.ccusage-sync.plist
    rm ~/Library/LaunchAgents/com.luoqi.ccusage-sync.plist

  改间隔：修改 StartInterval（秒）
  改路径：修改 ProgramArguments 里的 cc-usage 绝对路径
-->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.luoqi.ccusage-sync</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/luoqi/.local/bin/cc-usage</string>
        <string>sync</string>
    </array>

    <!-- 每 10 分钟跑一次 -->
    <key>StartInterval</key>
    <integer>300</integer>

    <!-- 加载时立即跑一次 -->
    <key>RunAtLoad</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/tmp/ccusage-sync.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/ccusage-sync.err</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
        <key>HOME</key>
        <string>/Users/luoqi</string>
    </dict>

    <!-- 错过后补跑（休眠唤醒后立即跑）-->
    <key>StartCalendarInterval</key>
    <dict>
        <key>Minute</key>
        <integer>0</integer>
    </dict>
</dict>
</plist>
